# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = service nginx start
user = root
autostart = true

[program:php5-fpm]
command = /usr/sbin/php5-fpm -FR
user = root
autostart = true
